# deploy/qbot-research.timer
[Unit]
Description=Run qbot research worker weekly

[Timer]
OnCalendar=weekly
Persistent=true
Unit=qbot-research.service

[Install]
WantedBy=timers.target
